# etc/supervisor.ini
[supervisord]
user = root
nodaemon=true

[program:pc]
directory = /app
command = pc run --env prod --no-frontend
stdout_logfile	= /var/log/applications/%(program_name)s.log
redirect_stderr = true

[program:nginx]
command = /usr/sbin/nginx
stdout_logfile	= /var/log/applications/%(program_name)s.log
redirect_stderr = true
